# manx-tts-user.service
[Unit]
Description=Gaelg AI Platform (User Service)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/exp/exp1/acp24csb/web_platform
Environment="PATH=/exp/exp1/acp24csb/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin"
Environment="WHISPER_HUB=/exp/exp1/acp24csb/model_instances/whisper/save"
Environment="WHISPER_DEVICE=cuda:1"
Environment="NLLB_DEVICE=cuda:0"
Environment="GRADTTS_ROOT=/exp/exp1/acp24csb/model_instances/Grad-TTS_graphemic"
Environment="GRAD_TTS_CKPT=/exp/exp1/acp24csb/model_instances/Grad-TTS_graphemic/checkpts/manx-22k.pt"
Environment="WHISPER_HUB=/exp/exp1/acp24csb/model_instances/whisper/save"
Environment="WHISPER_CKPT=/exp/exp1/acp24csb/model_instances/whisper/save/CKPT+2026-03-21+19-13-49+00"
Environment="NLLB_CHECKPOINTS=/exp/exp1/acp24csb/model_instances/nllb"
Environment="HF_HOME=/exp/exp1/acp24csb/hf_cache"
Environment="OUTPUT_DIR=/exp/exp1/acp24csb/web_platform/outputs"
Environment="LOG_DIR=/exp/exp1/acp24csb/web_platform/logs"
Environment="VC_ROOT=/exp/exp1/acp24csb/model_instances/voice_converter"
Environment="MALE_REF_DIR=/store/store1/data/HiFi-TTS/hi_fi_tts_v0/audio/9017_clean"
Environment="FEMALE_REF_DIR=/store/store1/data/LJSpeech-1.1/wavs"
Environment="OUTPUT_CLEANUP_TTL_HOURS=24"
Environment="OUTPUT_CLEANUP_INTERVAL_MINUTES=60"
Environment="GPU_MEMORY_WARNING_PERCENT=85"
Environment="SYNTHESIZE_TIMEOUT_SECONDS=30"
Environment="TRANSCRIBE_TIMEOUT_SECONDS=120"
Environment="TRANSLATE_TIMEOUT_SECONDS=30"
Environment="CONVERT_TIMEOUT_SECONDS=30"
Environment="MIN_DISK_FREE_MB=200"
Environment="RATE_LIMIT_REQUESTS=20"
Environment="RATE_LIMIT_WINDOW_SECONDS=60"
Environment="SKIP_TTS_WARMUP=false"
ExecStart=/exp/exp1/acp24csb/venv/bin/uvicorn backend.main:app --host 143.167.8.81 --port 8000
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=manx-tts

[Install]
WantedBy=default.target
